# install/eacdd.service
[Unit]
Description=eacd deployment daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/eacdd --config /etc/eacd/server.yaml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
